# deploy/quadlet/nutriscrape-flow.container
[Unit]
Description=NutriScrape batch pipeline as a parallel Prefect flow (run to completion)
After=neo4j.service
Requires=neo4j.service

[Container]
ContainerName=nutriscrape-flow
Image=localhost/nutriscrape:latest
Network=nectar.network
EnvironmentFile=%h/.config/nectar/platform.env
Environment=NEO4J_URI=bolt://neo4j:7687
Environment=NEO4J_USER=nutriscrape_writer
Secret=nutriscrape_pass,type=env,target=NEO4J_PASSWORD
Volume=corpus-staging:/data:Z
Volume=fdc-staging:/fdc:Z
# FDC_API_KEY is only the API fallback (unneeded once the bulk export is imported); uncomment to use:
# Secret=fdc_api_key,type=env,target=FDC_API_KEY
# Runs schema -> knowledge -> fdc-import -> (parallel ingest, NUTRISCRAPE_MAX_PARALLEL batches) ->
# cluster -> materialize as one Prefect DAG. Same image and env as the sequential nutriscrape unit.
Exec=flow

[Service]
Type=oneshot
# not enabled at boot; triggered manually or by an optional .timer
